# settings/systemd/system/mm8d.service.in
[Unit]
Description=MM8D operating software

[Service]
ExecStart=@prefix@/bin/mm8d.py
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target

