# settings/cron.d/mm8d-sw
# /etc/cron.d/mm8d-sw: crontab entries for MM8D

SHELL=/bin/sh
PATH=/usr/local/bin:/usr/bin

0 2 * * *			root	mm8d-maintainlog > /dev/null 2>&1
