# packaging/sw-armhf/prerm
#!/bin/sh
set -e
if [ \( "$1" = "upgrade" -o "$1" = "remove" \) -a -L /usr/doc/mm8d-sw ]; then
  rm -f /usr/doc/mm8d-sw
fi
invoke-rc.d mm8d.sh stop
systemctl disable mm8d.service
